# supervisord.conf
[supervisord]
nodaemon=true

[program:engineeringframework]
command=./MainModule

[program:cncprototype]
command=./CNCPrototype